# integrationTests/nmcontroller.service
[Unit]
After=network.target

[Service]
ExecStart=/opt/MockFog/NodeManager/nmController/runController.sh

[Install]
WantedBy=default.target